# conf/hades-deputy.service.in
[Unit]
Description=Hades deputy service for privileged operations
Documentation=https://agdsn.github.io/hades/
After=hades-auth-alternative-dns.service
After=hades-auth-pristine-dns.service
After=hades-database.service
After=hades-radius.service
After=hades-root-netns.service
BindsTo=hades-root-netns.service
Wants=hades-auth-alternative-dns.service
Wants=hades-auth-pristine-dns.service
Wants=hades-database.service
Wants=hades-radius.service

[Service]
Type=dbus
BusName=@DEPUTY_DBUS_NAME@
SyslogIdentifier=hades-deputy
Environment=PGAPPNAME=hades-deputy
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-deputy $DEPUTY_OPTS
Restart=always
RestartSec=5
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/deputy
;RuntimeDirectoryMode=0750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
