# conf/hades-database.service.in
[Unit]
Description=Hades database service (PostgreSQL)
Documentation=https://agdsn.github.io/hades/
After=hades-root-netns.service
BindsTo=hades-root-netns.service

[Service]
Type=forking
TimeoutSec=120
User=@DATABASE_USER@
Group=@DATABASE_GROUP@
SyslogIdentifier=hades-database
EnvironmentFile=@systemdenvfile@
ExecStartPre=@pkglibexecdir@/control-database.sh init
ExecStart=@pkglibexecdir@/control-database.sh start
ExecReload=@pkglibexecdir@/control-database.sh reload
ExecStop=@pkglibexecdir@/control-database.sh stop
Restart=always

[Install]
WantedBy=hades.target
