# conf/hades-auth-dhcp-leases.socket.in
[Socket]
ListenStream=@AUTH_DHCP_SCRIPT_SOCKET@
