# conf/hades-radius-vrrp.service.in
[Unit]
Description=Hades VRRP service for RADIUS services (keepalived)
Documentation=https://agdsn.github.io/hades/
After=hades-network.service hades-database.service
Wants=hades-network.service hades-database.service

[Service]
Type=forking
PIDFile=@pkgrunstatedir@/radius-vrrp/keepalived.pid
KillMode=process
SyslogIdentifier=hades-radius-vrrp
EnvironmentFile=@systemdenvfile@
ExecStartPre=@venvdir@/bin/hades-generate-config radius-keepalived.conf.j2 @pkgrunstatedir@/radius-vrrp/keepalived.conf
ExecStart=/usr/sbin/keepalived --log-console --vrrp --use-file=@pkgrunstatedir@/radius-vrrp/keepalived.conf --pid=@pkgrunstatedir@/radius-vrrp/keepalived.pid --vrrp_pid=@pkgrunstatedir@/radius-vrrp/keepalived-vrrp.pid
ExecReload=@venvdir@/bin/hades-generate-config radius-keepalived.conf.j2 @pkgrunstatedir@/radius-vrrp/keepalived.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
