# conf/hades-auth-pristine-dns.service.in
[Unit]
Description=Hades DNS service for authenticated users (unbound)
Documentation=https://agdsn.github.io/hades/
After=hades-network.service
Wants=hades-network.service

[Service]
Group=@AUTH_DNS_GROUP@
Type=simple
PIDFile=@AUTH_DNS_PRISTINE_PID_FILE@
SyslogIdentifier=hades-auth-pristine-dns
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config unbound-pristine.conf.j2 @pkgrunstatedir@/auth-dns/unbound-pristine.conf
ExecStartPre=@pkglibexecdir@/update-trust-anchor.sh
ExecStartPre=@UNBOUND_CHECKCONF@ @pkgrunstatedir@/auth-dns/unbound-pristine.conf
ExecStart=@IP@ netns exec @AUTH_NAMESPACE@ @UNBOUND@ -c @pkgrunstatedir@/auth-dns/unbound-pristine.conf $PRISTINE_UNBOUND_OPTS
ExecReload=@bindir@/hades-generate-config unbound-pristine.conf.j2 @pkgrunstatedir@/auth-dns/unbound-pristine.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
