# conf/hades-agent.service.in
[Unit]
Description=Hades agent (Celery Worker)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
Before=hades.service
After=hades-database.service
After=hades-root-netns.service
BindsTo=hades-database.service
BindsTo=hades-root-netns.service

[Service]
Type=simple
PIDFile=@pkgrunstatedir@/agent/worker.pid
KillMode=mixed
User=@AGENT_USER@
Group=@AGENT_GROUP@
SyslogIdentifier=hades-agent
Environment=PGAPPNAME=hades-agent
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-agent worker --pidfile=@pkgrunstatedir@/agent/worker.pid $AGENT_OPTS
Restart=always
RestartSec=5

;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/agent
;RuntimeDirectoryMode=2750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.service
