# conf/hades-agent.service.in
[Unit]
Description=Hades agent (Celery Worker)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
Before=hades.service
After=hades-database.service
After=hades-root-netns.service
BindsTo=hades-database.service
BindsTo=hades-root-netns.service

[Service]
Type=simple
PIDFile=@pkgrunstatedir@/agent/worker.pid
KillMode=mixed
User=@AGENT_USER@
Group=@AGENT_GROUP@
SyslogIdentifier=hades-agent
Environment=PGAPPNAME=hades-agent
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-agent --pid-file @pkgrunstatedir@/agent/agent.pid $AGENT_OPTS
Restart=always
RestartSec=5
# EX_USAGE, EX_UNAVAILABLE, EX_CONFIG
RestartPreventExitStatus=64 69 78

RuntimeDirectory=hades/agent
RuntimeDirectoryMode=2750
RuntimeDirectoryPreserve=yes
