# conf/hades-deputy.service.in
[Unit]
Description=Hades deputy service for privileged operations
Documentation=https://agdsn.github.io/hades/
After=hades-auth-dns.service hades-radius.service
Wants=hades-network.service

[Service]
Type=dbus
BusName=@DEPUTY_DBUS_NAME@
SyslogIdentifier=hades-deputy
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-deputy $DEPUTY_OPTS
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
