# conf/hades-cleanup.service.in
[Unit]
Description=Hades cleanup task (Service)
Documentation=https://agdsn.github.io/hades/
After=hades-database.service
Wants=hades-database.service

[Service]
Type=oneshot
User=@RADIUS_USER@
Group=@RADIUS_GROUP@
SyslogIdentifier=hades-cleanup
Environment=PGAPPNAME=hades-cleanup
ExecStart=@venvdir@/bin/hades-cleanup

[Install]
WantedBy=hades.target
