# conf/hades-unauth-http.service.in
[Unit]
Description=Hades http server (nginx)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service hades-unauth-portal.service
Wants=hades-network.service hades-unauth-portal.service

[Service]
Type=forking
PIDFile=@PORTAL_NGINX_PID_FILE@
SyslogIdentifier=hades-unauth-http
EnvironmentFile=@systemdenvfile@
ExecStartPre=@venvdir@/bin/hades-generate-config nginx @pkgrunstatedir@/unauth-http
ExecStart=@IP@ netns exec @UNAUTH_NAMESPACE@ nginx -c @pkgrunstatedir@/unauth-http/nginx.conf
ExecReload=@venvdir@/bin/hades-generate-config nginx @pkgrunstatedir@/unauth-http
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
KillSignal=SIGQUIT
KillMode=mixed
SyslogLevel=err
