# conf/hades-auth-dhcp.service.in
[Unit]
Description=Hades DHCP service for authenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service
Wants=hades-network.service

[Service]
Type=dbus
BusName=@AUTH_DHCP_DBUS_NAME@
SyslogIdentifier=hades-auth-dhcp
EnvironmentFile=@systemdenvfile@
ExecStartPre=@venvdir@/bin/hades-generate-config auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecStart=@IP@ netns exec @AUTH_NAMESPACE@ @DNSMASQ@ -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecReload=@venvdir@/bin/hades-generate-config auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
