# conf/hades-auth-vrrp.service.in
[Unit]
Description=Hades VRRP service for auth services (keepalived)
Documentation=https://agdsn.github.io/hades/
After=hades-auth-netns.service
BindsTo=hades-auth-netns.service
JoinsNamespaceOf=hades-auth-netns.service

[Service]
Type=dbus
BusName=@AUTH_VRRP_DBUS_NAME@
PrivateNetwork=yes
PIDFile=@pkgrunstatedir@/auth-vrrp/keepalived.pid
KillMode=process
SyslogIdentifier=hades-auth-vrrp
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config keepalived/auth.conf.j2 @pkgrunstatedir@/auth-vrrp/keepalived.conf
ExecStart=@KEEPALIVED@ --log-console --vrrp --use-file=@pkgrunstatedir@/auth-vrrp/keepalived.conf --pid=@pkgrunstatedir@/auth-vrrp/keepalived.pid --vrrp_pid=@pkgrunstatedir@/auth-vrrp/keepalived-vrrp.pid $AUTH_KEEPALIVED_OPTS
ExecReload=@bindir@/hades-generate-config keepalived/auth.conf.j2 @pkgrunstatedir@/auth-vrrp/keepalived.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/auth-vrrp
;RuntimeDirectoryMode=0750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
