# conf/hades-unauth-netns-cleanup.service.in
[Unit]
Description=Hades unauth network namespace cleanup
Documentation=https://bugzilla.redhat.com/show_bug.cgi?id=1353039

[Service]
Type=oneshot
ExecStart=-@IP@ netns delete unauth
