# conf/hades-unauth-portal.service.in
[Unit]
Description=Hades portal (Flask app on uWSGI)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
ReloadPropagatedFrom=hades.service
Before=hades.service
After=hades-database.service
After=hades-unauth-netns.service
BindsTo=hades-unauth-netns.service
JoinsNamespaceOf=hades-unauth-netns.service
Wants=hades-database.service

[Service]
Type=notify
User=@PORTAL_USER@
Group=@PORTAL_GROUP@
SyslogIdentifier=hades-unauth-portal
Environment=PGAPPNAME=hades-unauth-portal
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config uwsgi.ini.j2 @pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecStart=@UWSGI@ --ini=@pkgrunstatedir@/unauth-portal/uwsgi.ini $UWSGI_OPTS
ExecReload=@bindir@/hades-generate-config uwsgi.ini.j2 @pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
KillSignal=SIGINT
NotifyAccess=all

PrivateNetwork=yes

RuntimeDirectory=hades/unauth-portal
RuntimeDirectoryMode=2750
RuntimeDirectoryPreserve=yes
