# conf/hades-unauth-dns.service.in
[Unit]
Description=Hades DNS service for unauthenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
ReloadPropagatedFrom=hades.service
Before=hades.service
# The socket is strictly required for startup because dnsmasq needs to call
# `hades-dhcp-script init` initially
After=hades-unauth-dhcp-leases.socket
BindsTo=hades-unauth-dhcp-leases.socket
After=hades-unauth-netns.service
BindsTo=hades-unauth-netns.service
JoinsNamespaceOf=hades-unauth-netns.service

[Service]
Type=dbus
BusName=@UNAUTH_DNS_DBUS_NAME@
SyslogIdentifier=hades-unauth-dns
EnvironmentFile=@systemdenvfile@
Environment=HADES_DHCP_SCRIPT_SOCKET=@UNAUTH_DHCP_SCRIPT_SOCKET@
ExecStartPre=@bindir@/hades-generate-config --group=@UNAUTH_DNS_GROUP@ unauth-dnsmasq.conf.j2 @pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecStartPre=@DNSMASQ@ --test -k --enable-dbus=@UNAUTH_DNS_DBUS_NAME@ --conf-file=@pkgrunstatedir@/unauth-dns/dnsmasq.conf $UNAUTH_DNSMASQ_OPTS
ExecStart=@DNSMASQ@ -k --enable-dbus=@UNAUTH_DNS_DBUS_NAME@ --conf-file=@pkgrunstatedir@/unauth-dns/dnsmasq.conf $UNAUTH_DNSMASQ_OPTS
ExecReload=@bindir@/hades-generate-config --group=@UNAUTH_DNS_GROUP@ unauth-dnsmasq.conf.j2 @pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecReload=@DNSMASQ@ --test -k --enable-dbus=@UNAUTH_DNS_DBUS_NAME@ --conf-file=@pkgrunstatedir@/unauth-dns/dnsmasq.conf $UNAUTH_DNSMASQ_OPTS
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

PrivateNetwork=yes
# runtime directory controlled by tmpfile.d config

[Install]
Alias=hades-unauth-dhcp.service
