# conf/hades-forced-refresh.service.in
[Unit]
Description=Hades forced refresh task (Service)
Documentation=https://agdsn.github.io/hades/
# no PartOf=hades because triggered by timer (which has PartOf=hades)
After=hades-database.service
After=hades-deputy.service
Requires=hades-database.service
Requires=hades-deputy.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
User=@AGENT_USER@
Group=@AGENT_GROUP@
SyslogIdentifier=hades-refresh
ExecStart=@DBUS_SEND@ --system --type=method_call --print-reply --reply-timeout=30000 --dest=de.agdsn.hades.deputy /de/agdsn/hades/deputy de.agdsn.hades.deputy.Refresh boolean:true
