# conf/hades-agent.service.in
[Unit]
Description=Hades agent (Celery Worker)
Documentation=https://agdsn.github.io/hades/
After=hades-network.service hades-database.service
Wants=hades-network.service hades-database.service

[Service]
Type=simple
PIDFile=@pkgrunstatedir@/agent/worker.pid
KillMode=mixed
User=@AGENT_USER@
Group=@AGENT_GROUP@
SyslogIdentifier=hades-agent
Environment=PGAPPNAME=hades-agent
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-agent worker --pidfile=@pkgrunstatedir@/agent/worker.pid $AGENT_OPTS
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
