# conf/hades-auth-dhcp.service.in
[Unit]
Description=Hades DHCP service for authenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service
Wants=hades-network.service

[Service]
Type=dbus
BusName=de.agdsn.hades.auth_dnsmasq
SyslogIdentifier=hades-auth-dhcp
EnvironmentFile=@systemdenvdir@/hades
ExecStartPre=@venvdir@/bin/hades-generate-config auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecStart=@IP@ netns exec auth @DNSMASQ@ -k --enable-dbus=de.agdsn.hades.auth_dnsmasq --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecReload=@venvdir@/bin/hades-generate-config auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
