# conf/hades-unauth-portal.service.in
[Unit]
Description=Hades portal (Flask app on uWSGI)
Documentation=https://agdsn.github.io/hades/
After=hades-database.service
After=hades-unauth-netns.service
BindsTo=hades-unauth-netns.service
JoinsNamespaceOf=hades-unauth-netns.service
Wants=hades-database.service

[Service]
Type=notify
PrivateNetwork=yes
SyslogIdentifier=hades-unauth-portal
User=@PORTAL_USER@
Group=@PORTAL_GROUP@
EnvironmentFile=@systemdenvfile@
Environment=PGAPPNAME=hades-unauth-portal
ExecStartPre=@bindir@/hades-generate-config uwsgi.ini.j2 @pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecStart=@UWSGI@ --ini=@pkgrunstatedir@/unauth-portal/uwsgi.ini $UWSGI_OPTS
ExecReload=@bindir@/hades-generate-config uwsgi.ini.j2 @pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
KillSignal=SIGINT
StandardError=syslog
NotifyAccess=all
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/unauth-portal
;RuntimeDirectoryMode=2750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
