# conf/hades-auth-alternative-dns.service.in
[Unit]
Description=Hades alternative DNS service for authenticated users (unbound)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
ReloadPropagatedFrom=hades.service
Before=hades.service
After=hades-auth-netns.service
BindsTo=hades-auth-netns.service
JoinsNamespaceOf=hades-auth-netns.service

[Service]
Type=notify
PIDFile=@AUTH_DNS_ALTERNATIVE_PID_FILE@
Group=@AUTH_DNS_GROUP@
SyslogIdentifier=hades-auth-alternative-dns
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config unbound-alternative.conf.j2 @pkgrunstatedir@/auth-dns/unbound-alternative.conf
ExecStartPre=@pkglibexecdir@/update-trust-anchor.sh
ExecStartPre=@UNBOUND_CHECKCONF@ @pkgrunstatedir@/auth-dns/unbound-alternative.conf
ExecStart=@UNBOUND@ -c @pkgrunstatedir@/auth-dns/unbound-alternative.conf $ALTERNATIVE_UNBOUND_OPTS
ExecReload=@bindir@/hades-generate-config unbound-alternative.conf.j2 @pkgrunstatedir@/auth-dns/unbound-alternative.conf
ExecReload=@UNBOUND_CHECKCONF@ @pkgrunstatedir@/auth-dns/unbound-alternative.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

PrivateNetwork=yes
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/auth-dns
;RuntimeDirectoryMode=2750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.service
