# conf/hades-unauth-dns.service.in
[Unit]
Description=Hades DNS service for unauthenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
After=hades-network.service
Wants=hades-network.service

[Service]
Type=dbus
BusName=@UNAUTH_DNS_DBUS_NAME@
SyslogIdentifier=hades-unauth-dns
EnvironmentFile=@systemdenvfile@
ExecStartPre=@venvdir@/bin/hades-generate-config unauth-dnsmasq.conf.j2 @pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecStart=@IP@ netns exec @UNAUTH_NAMESPACE@ /usr/sbin/dnsmasq -k --enable-dbus=@UNAUTH_DNS_DBUS_NAME@ --conf-file=@pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecReload=@venvdir@/bin/hades-generate-config unauth-dnsmasq.conf.j2 @pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
Alias=hades-unauth-dhcp.service
