# conf/hades-root-netns.service.in
[Unit]
Description=Hades root network namespace
Documentation=https://github.com/systemd/systemd/issues/2741#issuecomment-336736214
StopWhenUnneeded=true

# Due to a bug in the systemd version shipped with Stretch ExecStopPost is not
# executed for Type=oneshot services. This has been fixed in an unknown, later
# version of systemd.
OnFailure=hades-root-netns-cleanup.service
Conflicts=hades-root-netns-cleanup.service

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStartPre=@bindir@/hades-generate-config netns/root/ @pkgrunstatedir@/root-netns/
ExecStartPre=@IFUP@ --interfaces=@pkgrunstatedir@/root-netns/interfaces --state-dir=@pkgrunstatedir@/root-netns --all --no-act

# Ask ip netns to create a named network namespace
# (This ensures that things like /var/run/netns are properly setup)
ExecStart=@IP@ netns add root

# Drop the network namespace that ip netns just created
ExecStart=@UMOUNT@ --verbose /var/run/netns/root

# Re-use the same name for the network namespace that systemd put us in, this
# should be root namespace.
ExecStart=@MOUNT@ --bind /proc/self/ns/net /var/run/netns/root
ExecStart=@SYSCTL@ --write net.ipv4.ip_nonlocal_bind=1 net.ipv4.ip_forward=1
ExecStart=@IFUP@ --verbose --no-loopback --interfaces=@pkgrunstatedir@/root-netns/interfaces --state-dir=@pkgrunstatedir@/root-netns --all
ExecStart=@IPTABLES_RESTORE@ @pkgrunstatedir@/root-netns/iptables-rules

# Teardown network interfaces
ExecStop=@IFDOWN@ --verbose --no-loopback --interfaces=@pkgrunstatedir@/root-netns/interfaces --state-dir=@pkgrunstatedir@/root-netns --all

# Delete the network namespace
ExecStopPost=-@IP@ netns delete root

;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/root-netns
;RuntimeDirectoryMode=0750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
