# conf/hades-radius.service.in
[Unit]
Description=Hades RADIUS service (FreeRADIUS)
Documentation=https://agdsn.github.io/hades/
After=hades-network.service hades-database.service
Wants=hades-network.service hades-database.service

[Service]
Type=forking
PIDFile=@RADIUS_PID_FILE@
SyslogIdentifier=hades-radius
Environment=PGAPPNAME=hades-radius
EnvironmentFile=@systemdenvfile@
ExecStartPre=@venvdir@/bin/hades-generate-config freeradius @pkgrunstatedir@/radius
ExecStartPre=@TOUCH@ @RADIUS_CLIENTS_FILE@
ExecStartPre=@RADIUSD@ -d @pkgrunstatedir@/radius -Cx -lstdout
ExecStart=@RADIUSD@ -d @pkgrunstatedir@/radius $FREERADIUS_OPTS
ExecReload=@venvdir@/bin/hades-generate-config freeradius @pkgrunstatedir@/radius
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
