# conf/hades-deputy.service.in
[Unit]
Description=Hades deputy service for privileged operations
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
ReloadPropagatedFrom=hades.service
Before=hades.service
After=hades-auth-alternative-dns.service
After=hades-auth-pristine-dns.service
After=hades-database.service
After=hades-radius.service
After=hades-root-netns.service
BindsTo=hades-root-netns.service
Wants=hades-auth-alternative-dns.service
Wants=hades-auth-pristine-dns.service
Wants=hades-database.service
Wants=hades-radius.service

[Service]
Type=dbus
BusName=@DEPUTY_DBUS_NAME@
SyslogIdentifier=hades-deputy
Environment=PGAPPNAME=hades-deputy
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-deputy --syslog $DEPUTY_OPTS
Restart=always
RestartSec=5
# EX_USAGE, EX_UNAVAILABLE, EX_CONFIG
RestartPreventExitStatus=64 69 78

RuntimeDirectory=hades/deputy
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes
