# conf/hades-unauth-http.service.in
[Unit]
Description=Hades http server (nginx)
Documentation=https://agdsn.github.io/hades/
After=hades-unauth-portal.service
After=hades-unauth-portal.service
BindsTo=hades-unauth-netns.service
JoinsNamespaceOf=hades-unauth-netns.service
Wants=hades-unauth-portal.service

[Service]
Type=forking
PrivateNetwork=yes
PIDFile=@PORTAL_NGINX_PID_FILE@
SyslogIdentifier=hades-unauth-http
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config --group=@PORTAL_GROUP@ nginx @pkgrunstatedir@/unauth-http
ExecStart=@NGINX@ -c @pkgrunstatedir@/unauth-http/nginx.conf $NGINX_OPTS
ExecReload=@bindir@/hades-generate-config --group=@PORTAL_GROUP@ nginx @pkgrunstatedir@/unauth-http
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
KillSignal=SIGQUIT
KillMode=mixed
SyslogLevel=err
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/unauth-http
;RuntimeDirectoryMode=2750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
