# conf/hades-auth-netns.service.in
[Unit]
Description=Hades auth network namespace
Documentation=https://github.com/systemd/systemd/issues/2741#issuecomment-336736214
RefuseManualStart=yes
StopWhenUnneeded=true
After=hades-root-netns.service
BindsTo=hades-root-netns.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/hades-generate-config netns/auth/ @pkgrunstatedir@/auth-netns/
ExecStart=@IFUP@ --interfaces=@pkgrunstatedir@/auth-netns/interfaces --state-dir=@pkgrunstatedir@/auth-netns --all --no-act

# Ask ip netns to create a named network namespace
# (This ensures that things like /var/run/netns are properly setup)
ExecStart=@IP@ netns add auth

# Drop the network namespace that ip netns just created
ExecStart=@UMOUNT@ --verbose /var/run/netns/auth

# Re-use the same name for the network namespace that systemd put us in
ExecStart=@MOUNT@ --verbose --bind /proc/self/ns/net /var/run/netns/auth

# Setup network settings
ExecStart=@SYSCTL@ --write net.ipv4.ip_nonlocal_bind=1 net.ipv4.ip_forward=1 net.ipv6.conf.all.disable_ipv6=1 net.ipv6.conf.default.disable_ipv6=1
ExecStart=@IFUP@ --verbose --interfaces=@pkgrunstatedir@/auth-netns/interfaces --state-dir=@pkgrunstatedir@/auth-netns --all
ExecStart=@IPSET@ -file @pkgrunstatedir@/auth-netns/ipsets restore
ExecStart=@IPTABLES_RESTORE@ @pkgrunstatedir@/auth-netns/iptables-rules

# Teardown network interfaces
ExecStop=@IFDOWN@ --verbose --interfaces=@pkgrunstatedir@/auth-netns/interfaces --state-dir=@pkgrunstatedir@/auth-netns --all

# Delete the network namespace
ExecStopPost=-@IP@ netns delete auth

PrivateNetwork=yes
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/auth-netns
;RuntimeDirectoryMode=0750
;RuntimeDirectoryPreserve=yes
