# conf/hades-unauth-dns.service.in
[Unit]
Description=Hades DNS service for unauthenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
After=hades-unauth-netns.service
BindsTo=hades-unauth-netns.service
JoinsNamespaceOf=hades-unauth-netns.service

[Service]
Type=dbus
PrivateNetwork=yes
BusName=@UNAUTH_DNS_DBUS_NAME@
SyslogIdentifier=hades-unauth-dns
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config unauth-dnsmasq.conf.j2 @pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecStart=@IP@ netns exec @UNAUTH_NAMESPACE@ @DNSMASQ@ -k --enable-dbus=@UNAUTH_DNS_DBUS_NAME@ --conf-file=@pkgrunstatedir@/unauth-dns/dnsmasq.conf $UNAUTH_DNSMASQ_OPTS
ExecReload=@bindir@/hades-generate-config unauth-dnsmasq.conf.j2 @pkgrunstatedir@/unauth-dns/dnsmasq.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/unauth-dns
;RuntimeDirectoryMode=2750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
Alias=hades-unauth-dhcp.service
