# conf/hades-unauth-portal.service.in
[Unit]
Description=Hades portal (Flask app on uWSGI)
Documentation=https://agdsn.github.io/hades/
After=hades-network.service hades-database.service
Wants=hades-network.service hades-database.service

[Service]
Type=notify
SyslogIdentifier=hades-unauth-portal
EnvironmentFile=@systemdenvfile@
Environment=PGAPPNAME=hades-unauth-portal
ExecStartPre=@venvdir@/bin/hades-generate-config uwsgi.ini.j2 @pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecStart=@IP@ netns exec @UNAUTH_NAMESPACE@ @UWSGI@ --ini=@pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecReload=@venvdir@/bin/hades-generate-config uwsgi.ini.j2 @pkgrunstatedir@/unauth-portal/uwsgi.ini
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
KillSignal=SIGINT
StandardError=syslog
NotifyAccess=all
