# conf/hades-auth-dhcp.service.in
[Unit]
Description=Hades DHCP service for authenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
After=hades-auth-netns.service
BindsTo=hades-auth-netns.service
JoinsNamespaceOf=hades-auth-netns.service

[Service]
Type=dbus
BusName=@AUTH_DHCP_DBUS_NAME@
SyslogIdentifier=hades-auth-dhcp
EnvironmentFile=@systemdenvfile@
ExecStartPre=@bindir@/hades-generate-config --group=@AUTH_DHCP_GROUP@ auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecStartPre=@DNSMASQ@ --test -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf $AUTH_DNSMASQ_OPTS
ExecStart=@DNSMASQ@ -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf $AUTH_DNSMASQ_OPTS
ExecReload=@bindir@/hades-generate-config --group=@AUTH_DHCP_GROUP@ auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecReload=@DNSMASQ@ --test -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf $AUTH_DNSMASQ_OPTS
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

PrivateNetwork=yes
;RuntimeDirectoryPreserve is not yet supported on Stretch's systemd
;RuntimeDirectory=hades/auth-dhcp
;RuntimeDirectoryMode=2750
;RuntimeDirectoryPreserve=yes

[Install]
WantedBy=hades.target
