# conf/hades-network.service.in
[Unit]
Description=Hades networking
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target

[Service]
Type=oneshot
RemainAfterExit=yes
SyslogIdentifier=hades-network
EnvironmentFile=@systemdenvfile@
ExecStart=@pkglibexecdir@/control-network.sh start
ExecStop=@pkglibexecdir@/control-network.sh stop

[Install]
WantedBy=multi-user.target
