# conf/hades-auth-dhcp-leases.service.in
[Unit]
Description=Hades auth DHCP lease management (accepting commands from `hades-dhcp-script`)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
ReloadPropagatedFrom=hades.service
# After= is implicitly added, however that doesn't ensure that thes unit is started
# exclusively by the socket (which needs to be accessible as an fd for the operation of this unit).
BindsTo=hades-auth-dhcp-leases.socket
After=hades-database.service
# db required for startup because we need a connection from the get-go
Requires=hades-database.service
# we could have bound to hades-auth-dhcp.service, but it might be useful
# to just be able to stop the dnsmasq and leave this server running,
# so it's more lenient to only propagate a stop on `hades.service`.
BindsTo=hades.service
JoinsNamespaceOf=hades-auth-netns.service

[Service]
Type=notify
SyslogIdentifier=hades-auth-dhcp-leases
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-lease-server
Restart=always
RestartSec=5
User=@AUTH_DHCP_USER@

PrivateNetwork=yes

RuntimeDirectory=hades/auth-dhcp
RuntimeDirectoryMode=2750
RuntimeDirectoryPreserve=yes
