# conf/hades-auth-dhcp.service.in
[Unit]
Description=Hades DHCP service for authenticated users (dnsmasq)
Documentation=https://agdsn.github.io/hades/
PartOf=hades.service
ReloadPropagatedFrom=hades.service
Before=hades.service
# The socket is strictly required for startup because dnsmasq needs to call
# `hades-dhcp-script init` initially
After=hades-auth-dhcp-leases.socket
BindsTo=hades-auth-dhcp-leases.socket
After=hades-auth-netns.service
BindsTo=hades-auth-netns.service
JoinsNamespaceOf=hades-auth-netns.service

[Service]
Type=dbus
BusName=@AUTH_DHCP_DBUS_NAME@
SyslogIdentifier=hades-auth-dhcp
EnvironmentFile=@systemdenvfile@
Environment=HADES_DHCP_SCRIPT_SOCKET=@AUTH_DHCP_SCRIPT_SOCKET@
ExecStartPre=@bindir@/hades-generate-config --group=@AUTH_DHCP_GROUP@ auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecStartPre=@DNSMASQ@ --test -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf $AUTH_DNSMASQ_OPTS
ExecStart=@DNSMASQ@ -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf $AUTH_DNSMASQ_OPTS
ExecReload=@bindir@/hades-generate-config --group=@AUTH_DHCP_GROUP@ auth-dnsmasq.conf.j2 @pkgrunstatedir@/auth-dhcp/dnsmasq.conf
ExecReload=@DNSMASQ@ --test -k --enable-dbus=@AUTH_DHCP_DBUS_NAME@ --conf-file=@pkgrunstatedir@/auth-dhcp/dnsmasq.conf $AUTH_DNSMASQ_OPTS
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5

PrivateNetwork=yes
# runtime directory controlled by tmpfile.d config
