# conf/hades-deputy.service.in
[Unit]
Description=Hades deputy service for privileged operations
Documentation=https://agdsn.github.io/hades/
After=hades-auth-dns.service hades-radius.service
After=hades-root-netns.service
BindsTo=hades-root-netns.service

[Service]
Type=dbus
BusName=@DEPUTY_DBUS_NAME@
SyslogIdentifier=hades-deputy
Environment=PGAPPNAME=hades-deputy
EnvironmentFile=@systemdenvfile@
ExecStart=@bindir@/hades-deputy $DEPUTY_OPTS
Restart=always
RestartSec=5

[Install]
WantedBy=hades.target
