# conf/hades-unauth-dhcp-leases.socket.in
[Socket]
ListenStream=@UNAUTH_DHCP_SCRIPT_SOCKET@
