# systemd/hades.target
[Unit]
Description=Hades
Documentation=https://agdsn.github.io/hades/

[Install]
WantedBy=multi-user.target
