# init/units/hades-radius-vrrp.service.in
[Unit]
Description=Hades VRRP service for RADIUS services (keepalived)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service hades-database.service
Wants=hades-network.service hades-database.service

[Service]
Type=forking
PIDFile=@PKGRUNSTATEDIR@/radius-vrrp/keepalived.pid
KillMode=process
SyslogIdentifier=hades-radius-vrrp
EnvironmentFile=@SYSTEMDENVDIR@/hades
ExecStartPre=@VENVDIR@/bin/hades-generate-config radius-keepalived.conf.j2 @PKGRUNSTATEDIR@/radius-vrrp/keepalived.conf
ExecStart=/usr/sbin/keepalived --log-console --vrrp --use-file=@PKGRUNSTATEDIR@/radius-vrrp/keepalived.conf --pid=@PKGRUNSTATEDIR@/radius-vrrp/keepalived.pid --vrrp_pid=@PKGRUNSTATEDIR@/radius-vrrp/keepalived-vrrp.pid
ExecReload=@VENVDIR@/bin/hades-generate-config radius-keepalived.conf.j2 @PKGRUNSTATEDIR@/radius-vrpp/keepalived.conf
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
