# init/units/hades-radius.service.in
[Unit]
Description=Hades RADIUS service (FreeRADIUS)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service hades-database.service
Wants=hades-network.service hades-database.service

[Service]
Type=forking
PIDFile=@pkgrunstatedir@/radius/radiusd.pid
SyslogIdentifier=hades-radius
Environment=PGAPPNAME=hades-radius
EnvironmentFile=@systemdenvdir@/hades
ExecStartPre=@venvdir@/bin/hades-generate-config freeradius @pkgrunstatedir@/radius
ExecStartPre=@RADIUSD@ -d @pkgrunstatedir@/radius -Cx -lstdout
ExecStart=@RADIUSD@ -d @pkgrunstatedir@/radius
ExecReload=@venvdir@/bin/hades-generate-config freeradius @pkgrunstatedir@/radius
ExecReload=@KILL@ -HUP $MAINPID
Restart=always
RestartSec=5
