# init/units/hades-database.service.in
[Unit]
Description=Hades database service (PostgreSQL)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service
Wants=hades-network.service

[Service]
Type=forking
TimeoutSec=120
SyslogIdentifier=hades-database
EnvironmentFile=@systemdenvdir@/hades
ExecStartPre=@pkglibexecdir@/control-database.sh init
ExecStart=@pkglibexecdir@/control-database.sh start
ExecReload=@pkglibexecdir@/control-database.sh reload
ExecStop=@pkglibexecdir@/control-database.sh stop
Restart=always
User=postgres
Group=postgres
