# init/units/hades-auth-dns.service.in
[Unit]
Description=Hades DNS service for authenticated users (unbound)
Documentation=https://agdsn.github.io/hades/
#PartOf=hades.target
After=hades-network.service
Wants=hades-network.service

[Service]
Group=@AUTH_DNS_GROUP@
Type=simple
PIDFile=@PKGRUNSTATEDIR@/auth-dns/unbound.pid
SyslogIdentifier=hades-auth-dns
EnvironmentFile=@SYSTEMDENVDIR@/hades
ExecStartPre=@VENVDIR@/bin/hades-generate-config unbound.conf.j2 @PKGRUNSTATEDIR@/auth-dns/unbound.conf
ExecStartPre=@PKGLIBEXECDIR@/update-trust-anchor.sh
ExecStartPre=@UNBOUND_CHECKCONF@ @PKGRUNSTATEDIR@/auth-dns/unbound.conf
ExecStart=@IP@ netns exec auth @UNBOUND@ -c @PKGRUNSTATEDIR@/auth-dns/unbound.conf
ExecReload=@VENVDIR@/bin/hades-generate-config unbound.conf.j2 @PKGRUNSTATEDIR@/auth-dns/unbound.conf
ExecReload=@UNBOUND@ -c @PKGRUNSTATEDIR@/auth-dns/unbound.conf reload
Restart=always
RestartSec=5
